# deploy/systemd/newdomofon-video-dvr.service
[Unit]
Description=NewDomofon Video DVR Engine
Documentation=file:/opt/newdomofon-video-node/docs/BAREMETAL_DEBIAN12.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=newdomofon
Group=newdomofon
WorkingDirectory=/opt/newdomofon-video-node/dvr-engine
EnvironmentFile=/etc/newdomofon-video/app.env
Environment=NODE_ENV=production
Environment=TZ=Europe/Moscow
ExecStart=/usr/bin/node /opt/newdomofon-video-node/dvr-engine/dist/index.js
Restart=always
RestartSec=5
TimeoutStopSec=45
KillSignal=SIGTERM
SyslogIdentifier=newdomofon-dvr
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/newdomofon-video /var/log/newdomofon-video /tmp
CapabilityBoundingSet=
AmbientCapabilities=
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
