# deploy/systemd/newdomofon-video-node-disk-guard.timer
[Unit]
Description=Run NewDomofon Video Node Disk Guard every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true
Unit=newdomofon-video-node-disk-guard.service

[Install]
WantedBy=timers.target
