# deploy/systemd/newdomofon-video-backend.service
[Unit]
Description=NewDomofon Video Backend API
Documentation=file:/opt/newdomofon-video/docs/BAREMETAL_DEBIAN12.md
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=newdomofon
Group=newdomofon
WorkingDirectory=/opt/newdomofon-video/backend
EnvironmentFile=/etc/newdomofon-video/app.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /opt/newdomofon-video/backend/dist/index.js
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
SyslogIdentifier=newdomofon-backend
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/log/newdomofon-video
CapabilityBoundingSet=
AmbientCapabilities=
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
