# deploy/systemd/newdomofon-smartyard-compat.service
[Unit]
Description=NewDomofon SmartYard Compatibility Proxy
Documentation=file:/opt/newdomofon-video-node/docs/BAREMETAL_DEBIAN12.md
After=network-online.target newdomofon-video-dvr.service
Wants=network-online.target

[Service]
Type=simple
User=newdomofon
Group=newdomofon
WorkingDirectory=/opt/newdomofon-video-node/smartyard-compat-proxy
EnvironmentFile=/etc/newdomofon-video/app.env
Environment=NODE_ENV=production
Environment=SMARTYARD_COMPAT_PORT=3082
Environment=DVR_ENGINE_URL=http://127.0.0.1:3010
ExecStart=/usr/bin/node /opt/newdomofon-video-node/smartyard-compat-proxy/server.js
Restart=always
RestartSec=3
TimeoutStopSec=20
KillSignal=SIGTERM
SyslogIdentifier=newdomofon-smartyard-compat
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/newdomofon-video /var/cache/newdomofon-video /var/log/newdomofon-video /tmp
StateDirectory=newdomofon-video
CacheDirectory=newdomofon-video
LogsDirectory=newdomofon-video
CapabilityBoundingSet=
AmbientCapabilities=
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
